# watchbox_0.9_1_all/DEBIAN/postinst
systemctl enable watchbox
systemctl start watchbox
